# scripts/launchd/com.kronos.auto-run.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kronos.auto-run</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/ethan/Kronos/scripts/run_mvp_auto_research.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ethan/Kronos</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/ethan/Kronos/logs/auto-run.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ethan/Kronos/logs/auto-run.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>KRONOS_SYNC_DATA</key>
    <string>1</string>
  </dict>
</dict>
</plist>
